# src/scraper/supervisord.conf
[supervisord]
nodaemon=true

[program:scraper]
command=python -m src.scraper.scraper_loop
directory=/app
autorestart=true
